# crontab
*/5 * * * * root freshclam >> /var/log/cron.log 2>&1
